# services/repeating_tasks/crontab
45 0 * * * python /app/main.py >> /var/log/cron.log 2>&1
# Don't forget to end this file with an empty newline
